# com.user.slack-presence.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.slack-presence</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Applications/SlackPresence.app/Contents/MacOS/SlackPresence</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/slack-presence.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/slack-presence.err</string>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>LowPriorityIO</key>
    <true/>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
